# puppet/modules/mastery_account/files/mastery_account.conf
start on runlevel [345]
stop on runlevel [0126]
expect fork
respawn
script
nodejs /usr/local/share/mastery/account/index.js /var/mastery/avatars/ 2>&1 | logger -t mastery_account &
end script
