# puppet/modules/mastery_server/files/mastery_server.conf
start on runlevel [345]
stop on runlevel [0126]
expect fork
respawn
script
nodejs /usr/local/share/mastery/server/lib/main.js 2>&1 | logger -t mastery_server &
end script

